# deploy/systemd/rdpproxy-port-watcher.path
[Unit]
Description=Watch .env for PUBLIC_PORT changes (RDPProxy)

[Path]
PathModified=/opt/rdpproxy/.env
Unit=rdpproxy-port-watcher.service

[Install]
WantedBy=multi-user.target
